# supervisord.conf
[supervisord]
nodaemon=true

[program:metube]
command=python3 -u /app/main.py
autorestart=true

[program:watcher]
command=/watcher.sh
autorestart=true
